# whitecoin

#! /bin/sh
# /etc/init.d/whitecoin

### BEGIN INIT INFO
# Provides:          whitecoin
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple script to start whitecoind at boot
# Description:       A simple script from www.stuffaboutcode.com which will start / stop whitecoind a boot / shutdown.
### END INIT INFO

# If you want a command to always run, put it here

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting whitecoind"
    # run application you want to start
    su - pi -c "/usr/local/bin/whitecoind -daemon"
    ;;
  stop)
    echo "Stopping whitecoind"
    # kill application you want to stop
    su - pi -c "/usr/local/bin/whitecoind stop"
    ;;
  *)
    echo "Usage: /etc/init.d/whitecoin {start|stop}"
    exit 1
    ;;
esac

exit 0 
